# Conf_Files/rtc_ds1307.conf
description "Sets up RTC and sets system time from RTC"

start on runlevel [2345]
stop on runlevel [!2345]

exec /usr/share/rtc-ds1307/clock_init.sh
